# env/config/ubuntu_appsotore.upstart.conf
description "start and stop the ubuntu appstore api"
version "1.0"
author "Brian Douglass"

start on runlevel [2345]
stop on runlevel [!2345]

exec /ubuntu-appstore/api/ubuntu_appstore.py
